# arm64v8/alpine/nginx/docker-alpine-nginx@.service
[Unit]
<<<<<<< HEAD:arm64v8/alpine/nginx/nginx@.service
Description=Containerized NGINX web server at LAN IP %i
=======
Description=Containerized Alpine NGINX web server at LAN IP %i
>>>>>>> 9b8f39f9fb24432e561ffdf28c217de5d79e53c1:arm64v8/alpine/nginx/docker-alpine-nginx@.service

Requires=docker.service

After=docker.service

[Service]
Type=simple

TimeoutStartSec=0

KillMode=none

EnvironmentFile=/etc/environment

ExecStartPre=/usr/bin/docker create \
    --name nginx \
    -e NGINX_DOMAIN_NAME="${NGINX_DOMAIN_NAME}" \
    -e NGINX_UPSTREAM_SERVERS="${NGINX_UPSTREAM_SERVERS}" \
    -p %i:80:80/tcp \
    -p %i:8000:8000/tcp \
    yuanmodu/arm64v8-alpine-nginx
ExecStart=/usr/bin/docker start -a nginx

ExecStop=/usr/bin/docker stop nginx
ExecStopPost=/usr/bin/docker rm nginx

[Install]
WantedBy=multi-user.target
